# zinuit/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:{{ zinuit_name }}-metel-web]
command={{ zinuit_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} -t {{ http_timeout }} metel.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/web.log
stderr_logfile={{ zinuit_dir }}/logs/web.error.log
user={{ user }}
directory={{ sites_dir }}

{% if use_rq %}
[program:{{ zinuit_name }}-metel-schedule]
command={{ zinuit_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/schedule.log
stderr_logfile={{ zinuit_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ zinuit_dir }}

[program:{{ zinuit_name }}-metel-default-worker]
command={{ zinuit_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/worker.log
stderr_logfile={{ zinuit_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ zinuit_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ zinuit_name }}-metel-short-worker]
command={{ zinuit_cmd }} worker --queue short
priority=4
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/worker.log
stderr_logfile={{ zinuit_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ zinuit_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ zinuit_name }}-metel-long-worker]
command={{ zinuit_cmd }} worker --queue long
priority=4
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/worker.log
stderr_logfile={{ zinuit_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ zinuit_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

{% else %}
[program:{{ zinuit_name }}-metel-workerbeat]
command={{ zinuit_dir }}/env/bin/python -m metel.celery_app beat -s beat.schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/workerbeat.log
stderr_logfile={{ zinuit_dir }}/logs/workerbeat.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ zinuit_name }}-metel-worker]
command={{ zinuit_dir }}/env/bin/python -m metel.celery_app worker -n jobs@%%h -Ofair --soft-time-limit 360 --time-limit 390 --loglevel INFO
priority=4
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/worker.log
stderr_logfile={{ zinuit_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=400
directory={{ sites_dir }}
killasgroup=true

[program:{{ zinuit_name }}-metel-longjob-worker]
command={{ zinuit_dir }}/env/bin/python -m metel.celery_app worker -n longjobs@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/worker.log
stderr_logfile={{ zinuit_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

[program:{{ zinuit_name }}-metel-async-worker]
command={{ zinuit_dir }}/env/bin/python -m metel.celery_app worker -n async@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/worker.log
stderr_logfile={{ zinuit_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

{% endif %}

[program:{{ zinuit_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/redis-cache.log
stderr_logfile={{ zinuit_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ zinuit_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/redis-queue.log
stderr_logfile={{ zinuit_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}

{% if metel_version > 1 %}
[program:{{ zinuit_name }}-redis-socketio]
command={{ redis_server }} {{ redis_socketio_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/redis-socketio.log
stderr_logfile={{ zinuit_dir }}/logs/redis-socketio.error.log
user={{ user }}
directory={{ sites_dir }}

{% if node %}
[program:{{ zinuit_name }}-node-socketio]
command={{ node }} {{ zinuit_dir }}/apps/metel/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ zinuit_dir }}/logs/node-socketio.log
stderr_logfile={{ zinuit_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ zinuit_dir }}
{% endif %}

{% endif %}

[group:{{ zinuit_name }}-web]
programs={{ zinuit_name }}-metel-web {%- if node -%} ,{{ zinuit_name }}-node-socketio {%- endif%}

{% if use_rq %}

[group:{{ zinuit_name }}-workers]
programs={{ zinuit_name }}-metel-schedule,{{ zinuit_name }}-metel-default-worker,{{ zinuit_name }}-metel-short-worker,{{ zinuit_name }}-metel-long-worker

{% else %}

[group:{{ zinuit_name }}-workers]
programs={{ zinuit_name }}-metel-workerbeat,{{ zinuit_name }}-metel-worker,{{ zinuit_name }}-metel-longjob-worker,{{ zinuit_name }}-metel-async-worker

{% endif %}

[group:{{ zinuit_name }}-redis]
programs={{ zinuit_name }}-redis-cache,{{ zinuit_name }}-redis-queue {%- if metel_version > 1 -%} ,{{ zinuit_name }}-redis-socketio {%- endif %}
